# container.target
[Unit]
Description=Special target for running systemd in docker containers
